# data/systemd/laniakea-spark.service.in
[Unit]
Description=Laniakea Spark
After=syslog.target network.target
ConditionPathExists=/etc/laniakea/spark.json

[Service]
Type=simple
Restart=on-failure
User=lkspark
Group=lkspark
ExecStart=@LKSPARK_INSTALL_BIN@

[Install]
WantedBy=multi-user.target
